# recipes/thestage/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0

[program:models]
directory=/opt/project
command=sh -c 'elastic-models serve diffusion "${MODEL_REPO}" --size "${MODEL_SIZE}" --batch "${MODEL_BATCH_SIZE}" --enable-metadata-server'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3

[program:nginx]
command=sh -c 'AUTH_TOKEN="${AUTH_TOKEN}" /usr/local/openresty/nginx/sbin/nginx -g "daemon off;"'
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0